# etc/systemd/system/napd.service
[Unit]
Description=napd auto-sleep daemon

[Service]
ExecStart=/usr/bin/napd

[Install]
WantedBy=multi-user.target
